# templates/default/systemd/chef-client.service.erb
[Unit]
Description=Chef Client daemon
After=network.target auditd.service

[Service]
Type=<%= @type %>
EnvironmentFile=<%= @sysconfig_file %>
ExecStart=<%= @client_bin %> <%= @exec_options %>
ExecReload=/bin/kill -HUP $MAINPID
SuccessExitStatus=3
<% if @restart_mode %>
Restart=<%= @restart_mode %>
<% end %>

[Install]
WantedBy=multi-user.target
