# templates/default/fedora/systemd/chef-client.service.erb
[Unit]
Description=Chef Client daemon
After=network.target auditd.service

[Service]
EnvironmentFile=<%= @sysconfig_file %>
ExecStart=<%= @client_bin %> -c $CONFIG -i $INTERVAL -s $SPLAY
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
