# templates/freebsd/chef-client.erb
#!/bin/sh

# PROVIDE: chef
# REQUIRE: LOGIN
# KEYWORD: nojail shutdown

. /etc/rc.subr

name="chef"
pidfile="<%= node["chef_client"]["run_path"] %>/${name}.pid"
command="<%= @client_bin %>"
command_interpreter="<%= RbConfig.ruby %>"
command_args="-i <%= node["chef_client"]["interval"] %> -s <%= node["chef_client"]["splay"] %> -d -L <%= node["chef_client"]["log_dir"] %>/client.log -P <%= node["chef_client"]["run_path"] %>/${name}.pid"
load_rc_config $name
run_rc_command "$1"
