# templates/default/arch/systemd/chef-client.service.erb
[Unit]
Description=Chef Client daemon
After=syslog.target network.target auditd.service

[Service]
Type=forking
EnvironmentFile=<%= @environment_file %>
ExecStart=<%= @client_bin %> -d -c $CONFIG -P $PIDFILE -i $INTERVAL -s $SPLAY
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=<%= node["chef_client"]["run_path"]%>/client.pid

[Install]
WantedBy=multi-user.target
